# src/vhost/build/modules/build/files/etc/supervisor/conf.d/nginx_stdout.conf
[program:nginx_stdout]
command = /usr/bin/tail -F /var/log/nginx/stdout.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
