# src/vhost/build/modules/vhost/files/etc/supervisor/conf.d/dockergen.conf
[program:dockergen-1]
command = /usr/local/bin/docker-gen -watch -only-exposed -notify "/usr/sbin/nginx -s reload" /root/vhost.conf.tmpl /etc/nginx/conf.d/vhost.conf
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[program:dockergen-2]
command = /usr/local/bin/docker-gen -watch -only-exposed -notify "/bin/chown -R nginx.nginx /var/www" /root/vhost.php.tmpl /var/www/vhost.php
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
