# src/build/modules/vhost/files/etc/supervisor/conf.d/dns.conf
[program:dns]
command = /usr/bin/node /opt/dns/dns.js
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
