# src/build/modules/vhost/files/etc/supervisor/conf.d/dockergen.conf
[program:dockergen-1]
command = /usr/local/bin/docker-gen -watch -only-exposed -notify "/usr/sbin/nginx -s reload" /root/vhost.conf.tmpl /etc/nginx/conf.d/vhost.conf
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true

[program:dockergen-2]
command = /usr/local/bin/docker-gen -watch -only-exposed -notify "/bin/chown -R nginx.nginx /vhost/data" /root/vhost.php.tmpl /vhost/data/vhost.php
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
